# mac/com.me.start-redshift.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.em.start-redshift</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/redshift</string>
  </array>

  <key>Nice</key>
  <integer>3</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/start-redshift.err</string>

  <key>StandardOutPath</key>
  <string>/tmp/start-redshift.out</string>
</dict>
</plist>
